# systemd/drs-api.service
[Unit]
Description=DRS Unterrichtsmaterial API (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=drs
Group=drs
WorkingDirectory=/opt/drs/app
EnvironmentFile=/etc/drs/config.env
EnvironmentFile=-/etc/drs/onlyoffice.env
ExecStart=/opt/drs/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/drs/data /opt/drs/app/data

[Install]
WantedBy=multi-user.target
